# guest/rootfs/alpine/executa.sh
#!/sbin/openrc-run
  
command="/usr/bin/node"
command_args="node_modules/@stencila/executa/dist/lib/cli.js --stdio"
