# guest/rootfs/alpine/execution-server.sh
#!/sbin/openrc-run
  
command="/usr/bin/node"
command_args="/usr/bin/execution-server"
