# guest/rootfs/alpine/hostname.sh
#!/sbin/openrc-run
  
command="hostname"
command_args="sparkla-alpine-vm"
